# supervisord_snmp.conf
[program:snmptrapd]
command=/usr/sbin/snmptrapd -CdfLf trap.log --disableAuthorization=yes
redirect_stderr=true
